# AbletonPerformanceHelper/AbletonPerformanceHelper/InstallerPkg/Payload/Library/LaunchDaemons/com.maas1ve.APH.Helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.maas1ve.APH.Helper</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/PrivilegedHelperTools/com.maas1ve.AbletonPerformanceHelper.helper</string>
  </array>
  <key>MachServices</key>
  <dict>
    <key>com.maas1ve.APH.Helper</key><true/>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/var/log/APHHelper.out.log</string>
  <key>StandardErrorPath</key><string>/var/log/APHHelper.err.log</string>
</dict>
</plist>
